# scripts/launchd/com.rjm.master.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rjm.master</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/motomoto/bin/rjm-run-agent.sh</string>
    <string>master</string>
  </array>

  <!-- Every 3 hours = 8×/day -->
  <key>StartInterval</key>
  <integer>10800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/logs/launchd_master.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/logs/launchd_master.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/motomoto</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
