# scripts/launchd/com.rjm.youtube-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rjm.youtube-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/motomoto/bin/rjm-youtube-review-daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/motomoto/Library/Logs/rjm-youtube-review-daily.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/motomoto/Library/Logs/rjm-youtube-review-daily.err.log</string>
</dict>
</plist>
