# scripts/launchd/com.rjm.research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rjm.research</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/motomoto/bin/rjm-run-agent.sh</string>
    <string>research</string>
  </array>

  <!-- Every 4 hours = 6×/day (offset from discover by 2 hours) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>17</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/logs/launchd_research.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/logs/launchd_research.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/motomoto</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
