# scripts/launchd/com.rjm.viral-learning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rjm.viral-learning</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/rjm.py</string>
    <string>content</string>
    <string>learning</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre</string>

  <!-- Daily at 18:00 local time (CET) — Learning Loop -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>18</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/logs/launchd_viral_learning.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/motomoto/Documents/Robert-Jan Mastenbroek Command Centre/logs/launchd_viral_learning.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/motomoto</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
